# driver-kafka/deploy/ssd-deployment/templates/gateway.service
[Unit]
Description=Conduktor Gateway
After=network.target

[Service]
ExecStart=/bin/java -jar /opt/gateway/gateway.jar
WorkingDirectory=/opt/gateway
RestartSec=1s
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
